# compose.service
[Unit]
Description=start docker compose containers
After=network.target

[Service]
Type=simple
User=pi
Group=pi
Environment="SCRIPTS=/home/pi/Documents/shellScripts"
Environment="REPO_NAME=zpwr"
Environment="ZSH_COMP_REPO_NAME=zsh-more-completions"
ExecStart=/usr/bin/env bash /home/pi/Documents/shellScripts/startContainers.sh
ExecStop=/usr/bin/env kill $MAINPID

[Install]
WantedBy=multi-user.target
