# compose.service
[Unit]
Description=start docker compose containers
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/pi/docker
ExecStart=/usr/bin/env docker-compose -f /home/pi/docker/zabbix-docker/docker-compose_v3_ubuntu_pgsql_latest.yaml -f /home/pi/docker/docker-elk/docker-compose.yml up -d
ExecStop=/usr/bin/env docker-compose -f /home/pi/docker/zabbix-docker/docker-compose_v3_ubuntu_pgsql_latest.yaml -f /home/pi/docker/docker-elk/docker-compose.yml down

[Install]
WantedBy=multi-user.target
